# deploy/gunicorn.service.example
# Gunicorn systemd 服务配置示例
# 复制为 gunicorn.service 并根据实际情况修改

[Unit]
Description=Gunicorn service for Django project
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/path/to/project
ExecStart=/path/to/project/venv/bin/gunicorn --access-logfile - --workers 3 --bind 127.0.0.1:8000 config.wsgi:application
Restart=on-failure

[Install]
WantedBy=multi-user.target