# deploy/supervisor.conf.example
# Supervisor 配置示例
# 复制为 supervisor.conf 并根据实际情况修改

[program:django_app]
directory=/path/to/project
command=/path/to/project/venv/bin/gunicorn --access-logfile - --workers 3 --bind 127.0.0.1:8000 config.wsgi:application
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/django_app.log
